# Win10/config/waybar/modules/updates/updates.timer
[Unit]
Description=Check for Pacman/AUR/git package updates

[Timer]
OnStartupSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target